# deploy/kalshi-collector.service
[Unit]
Description=Kalshi Market Data Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kalshi
WorkingDirectory=/opt/kalshi-bot
EnvironmentFile=/opt/kalshi-bot/.env
ExecStart=/opt/kalshi-bot/.venv/bin/python -m kalshi_bot.data_collector \
    --db /var/lib/kalshi/market_data.db \
    --interval 60 \
    --log-level INFO
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kalshi-collector

# Safety: if it crashes more than 5 times in 10 minutes, stop trying
StartLimitIntervalSec=600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
